# contrib/init/atlasd.conf
description "Atlas Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env ATLASD_BIN="/usr/bin/atlasd"
env ATLASD_USER="atlas"
env ATLASD_GROUP="atlas"
env ATLASD_PIDDIR="/var/run/atlasd"
# upstart can't handle variables constructed with other variables
env ATLASD_PIDFILE="/var/run/atlasd/atlasd.pid"
env ATLASD_CONFIGFILE="/etc/atlas/atlas.conf"
env ATLASD_DATADIR="/var/lib/atlasd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # atlasd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$ATLASD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run atlasd."
        echo "The setting must appear in $ATLASD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Atlas Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$ATLASD_PIDDIR"
    chmod 0755 "$ATLASD_PIDDIR"
    chown $ATLASD_USER:$ATLASD_GROUP "$ATLASD_PIDDIR"
    chown $ATLASD_USER:$ATLASD_GROUP "$ATLASD_CONFIGFILE"
    chmod 0660 "$ATLASD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$ATLASD_PIDFILE" \
    --chuid $ATLASD_USER:$ATLASD_GROUP \
    --exec "$ATLASD_BIN" \
    -- \
    -pid="$ATLASD_PIDFILE" \
    -conf="$ATLASD_CONFIGFILE" \
    -datadir="$ATLASD_DATADIR" \
    -disablewallet \
    -daemon

